# files/1password-flatpak-setup.service
[Unit]
Description=Wire Flatpak browsers up to the native 1Password app
Documentation=https://github.com/reinier/rheniite
After=graphical-session.target
ConditionPathExists=/usr/libexec/1Password-BrowserSupport

[Service]
Type=oneshot
ExecStart=/usr/libexec/rheniite/1password-flatpak-setup
# A failure here must never block the user session.
SuccessExitStatus=0 1

[Install]
WantedBy=default.target
